# desktop/airootfs/etc/systemd/system/getty@tty2.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin bob --noclear %I 38400 linux
