# serious_face/airootfs/usr/lib/systemd/system/rtorrent.service
[Unit]
Description=tmux

[Service]
User=user
ExecStart=/usr/bin/tmux new-session -s tmux_rtorrent -n rtorrent -d rtorrent



